# packages/shells-bridge/docker/minion-vnc.service
[Unit]
Description=Minion loopback VNC server
After=minion-x-display.service
Requires=minion-x-display.service

[Service]
Type=simple
User=agent
Group=agent
Environment=DISPLAY=:1
ExecStart=/usr/bin/x11vnc -display :1 -localhost -forever -shared -nopw -rfbport 5900
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
